# born2beroot/vtian/.config/systemd/user/monitoring.service
[Unit]
Description=Run monitoring script

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitoring.sh >> /var/log/monitoring.log 2>&1